# roles/sys/sysstat/templates/sysstat-summary.timer.j2
[Unit]
Description=Generate summary of yesterday's process accounting

[Timer]
OnCalendar=00:07:00

[Install]
WantedBy=sysstat.service
