# roles/sys/sysstat/templates/sysstat-collect.service.j2
[Unit]
Description=system activity accounting tool
Documentation=man:sa1(8)
After=sysstat.service

[Service]
Type=oneshot
User=root
ExecStart={{ sysstat_sa1_bin }} -S XALL 60 1
