# roles/sys/sysstat/templates/sysstat-collect.timer.j2
[Unit]
Description=Run system activity accounting tool every minute

[Timer]
OnCalendar=*:00/1

[Install]
WantedBy=sysstat.service
