# roles/postgres/initdb/templates/postgresql.service.j2
.include /lib/systemd/system/postgresql-{{pgversion}}.service

[Service]
Environment=PGDATA={{ postgres_data_dir }}
